# scripts/auto-cpufreq.service
[Unit]
Description=auto-cpufreq - Automatic CPU speed & power optimizer for Linux
After=network.target network-online.target
ConditionPathExists=/var/log/auto-cpufreq.log

[Service]
Type=simple
User=root
ExecStart=/usr/bin/sh -c '/usr/bin/python3 /usr/bin/auto-cpufreq --daemon'
StandardOutput=append:/var/log/auto-cpufreq.log
[Install]
WantedBy=multi-user.target
